# openmrs/gateway/crontab
# m h  dom mon dow   command
0 5  * * *  ../cron_job.sh
